# ping-base/ping.upstart
description "ping"

start on runlevel [2345]
stop on runlevel [016]
setuid vagrant
respawn

exec /var/local/ping/run.sh